# files/acmetool.crontab
0 0 * * * /usr/bin/acmetool --response-file /var/lib/acme/conf/responses --batch
